# applications/kdapps/kaspa-auth/systemd/kaspa-auth.service
[Unit]
Description=Kaspa Auth Daemon (user)
After=default.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
## Override this path if you install the binary elsewhere
## Default: run in dev mode (file-backed wallet under ~/.local/share/kaspa-auth)
ExecStart=%h/.cargo/bin/kaspa-auth --dev-mode daemon start --foreground \
  --socket-path %t/kaspa-auth.sock \
  --data-dir %h/.local/share/kaspa-auth \
  --session-timeout 3600
Environment=KASPA_AUTH_DATA_DIR=%h/.local/share/kaspa-auth
Restart=always
RestartSec=2s

# Hardening
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=%h/.local/share %t
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=default.target
